# meta-facebook/meta-yosemite4/recipes-phosphor/mctp/files/setup-local-eid.service
[Unit]
Description=MCTP local network configuration
Before=mctpd.service

[Service]
Type=oneshot
ExecStart=/usr/share/mctp/mctp-config.sh

[Install]
WantedBy=multi-user.target
