# meta-facebook/meta-catalina/recipes-catalina/hmc-util/files/sync-hmc-datetime.timer
[Unit]
Description=Run sync-hmc-datetime service every hour 

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target